# rc.cloud/files/etc/init.d/rc.cloud-final
#!/bin/sh /etc/rc.common
#  Финальный хук – просто отметка в логах
START=92
USE_PROCD=1

start_service() {
	logger -t rc.cloud-final "===== rc.cloud-final: system completely up ====="
}
stop_service() { :; }
