# rc.cloud/files/etc/init.d/rc.cloud-final
#!/bin/sh /etc/rc.common
# cloud-init lite – final phase (marker)
START=92
STOP=10
USE_PROCD=1

start_service() {
	logger -t rc.cloud-final "===== rc.cloud-final: system completely up at $(date) ====="
}
stop_service() { :; }
