# rc.cloud/files/etc/init.d/rc.cloud-final
#!/bin/sh /etc/rc.common
# Финальная стадия: просто лог + отмонтирование ISO

START=92
STOP=10
USE_PROCD=1

CIDATA_DIR="/mnt/cidata"

start_service() {
	logger -t rc.cloud-final "===== rc.cloud-final: system completely up ====="

	if mountpoint -q "$CIDATA_DIR"; then
		umount "$CIDATA_DIR"
		rmdir  "$CIDATA_DIR"
	fi
}

stop_service() { :; }
